# src/keepalive.service
[Unit]
Description=Restart the system after power outages

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/keepalive.sh
