# src/wmp.service
[Unit]
Description=WakeMyPotato: Restart the system after power outages

[Service]
Type=oneshot
ExecStart=/opt/wmp/wmp-run 600 y
