# src/wmp.service
[Unit]
Description=WakeMyPotato: Restart the system after power outages

[Service]
Type=oneshot
# ExecStart=/usr/local/sbin/wmp.sh 600
